# deploy/netboxvmsync.timer
[Unit]
Description=Run netboxvmsync daily

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target